# files/2_monitor/supervisor/ftlab/DSM101.conf
[program:run_DSM101]
command = /opt/monitor/script/DSM101.py
process_name = %(program_name)s
autostart = false
autorestart = true
user = pi
redirect_stderr = true
stdout_logfile = /opt/monitor/log/DSM101.out
stderr_logfile = /opt/monitor/log/DSM101.err
environment =
	    DEVICE_USB = "/dev/serial/by-id/usb-FTDI_FT232R_USB_UART_A50285BI-if00-port0",
	    DEVICE_POS = "XXXXX",
	    DB_URL     = "https://influxdb.monitor.yemilab.kr",
	    DB_ORG     = "Yemilab",
	    DB_BUCKET  = "ENVMONITOR_BUCKET",
	    DB_TOKEN   = "XXXX"
	    