# files/2_monitor/supervisor/ADAM_6018p.conf
[program:run_adam6018p]
command = /opt/monitor/script/ADAM_6018p.py
process_name = %(program_name)s
autostart = false
autorestart = true
user = pi
redirect_stderr = true
stdout_logfile = /opt/monitor/log/ADAM_6018p.out
stderr_logfile = /opt/monitor/log/ADAM_6018p.err
environment =
	DEVICE_POS = "XXXX",
	DEVICE_IP  = "61",
	DEVICE_NCH = '8',
	DB_URL     = "https://influxdb.monitor.yemilab.kr",
	DB_ORG     = "Yemilab",
	DB_BUCKET  = "ENVMONITOR_BUCKET",
	DB_TOKEN   = "XXX"
