# files/2_monitor/supervisor/PCInfo.conf
[program:run_PCInfo]
command = /opt/monitor/script/PCInfo.py
process_name = %(program_name)s
autostart = false
autorestart = true
user = pi
redirect_stderr = true
stdout_logfile = /opt/monitor/log/PCInfo.out
stderr_logfile = /opt/monitor/log/PCInfo.err
environment =
	DEVICE_DEV = "Monitor",
	DEVICE_POS = "XXXXXX",
	DB_URL     = "https://influxdb.monitor.yemilab.kr",
	DB_ORG     = "Yemilab",
	DB_BUCKET  = "ENVMONITOR_BUCKET",
	DB_TOKEN   = "XXXX"
